# bin/zram-service/zram-swap.service
[Unit]
Description=ZRAM swap setup
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/zram-setup.sh start
ExecStop=/usr/local/bin/zram-setup.sh stop

[Install]
WantedBy=multi-user.target
